# .systemd/cepex.service
[Unit]
Description=CepexApp

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=%h/cepex-server/app

[Install]
WantedBy=default.target
